# deploy/systemd/ai-calendar-adviser-pull-deploy.service
[Unit]
Description=ai-calendar-adviser: pull main + redeploy if HEAD changed
Documentation=https://github.com/hamlet-archer/ai-calendar-adviser
After=network-online.target ai-calendar-adviser.service
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
ExecStart=/opt/ai-calendar-adviser-deploy/pull-deploy.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ai-calendar-adviser-pull-deploy
TimeoutStartSec=600

# Hardening. NoNewPrivileges intentionally NOT set — pull-deploy needs
# `sudo systemctl restart …` for the ai-calendar-adviser service + sync
# timer. ProtectHome relaxed to read-only so git can read the staged
# deploy key + known_hosts (matches the ai-comms-adviser pull-deploy
# pattern, which itself patterns ai-chief-pull-deploy).
ProtectSystem=full
# Carve a write hole for /etc/systemd/system so deploy.sh's `cp` of
# updated unit files succeeds — without this, ProtectSystem=full made
# the cp fail with EROFS after the script had already stopped the
# services (the ai-comms-adviser outage on 2026-05-10 17:31 UTC
# documents this exact failure mode).
ReadWritePaths=/etc/systemd/system
ProtectHome=read-only
PrivateTmp=true

[Install]
WantedBy=multi-user.target
