# deploy/systemd/ai-calendar-adviser.service
[Unit]
Description=AI Calendar Adviser — calendar.{query,find_free_slot}.v1 Unix-socket RPC server
Documentation=https://github.com/hamlet-archer/ai-calendar-adviser
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
# Type=simple rather than Type=notify: Node 20 LTS does not expose
# AF_UNIX dgram sockets, so a stock-Node sd_notify(READY=1) is not
# implementable without a native addon. Pragmatic mirror of
# ai-comms-adviser-query.service (also Type=simple). Readiness is
# confirmed by the ExecStartPost probe below.
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/ai-calendar-adviser
EnvironmentFile=/opt/ai-calendar-adviser/.env
ExecStart=/usr/bin/node dist/main.js
Restart=on-failure
# Slow restart so a Google quota rate-limit doesn't re-burn the per-minute
# window on a crash-loop (folded from the timeout-override.conf drop-in).
RestartSec=120
StandardOutput=journal
StandardError=journal
# Cold boot self-check + initial sync over the calendars + quota-retry
# headroom. A bare 30s was insufficient: the boot check is ~5s but Google
# API quota retries can stretch it to 60-90s under sustained rate-limit, so
# 300s gives the googleapis exponential-backoff path room (folded from the
# timeout-override.conf drop-in so the running unit == this file; B8.10.3).
TimeoutStartSec=300
TimeoutStopSec=10

# Block until the RPC socket is actually bound — the unit isn't "active"
# until callers can connect. Caps the boot race ai-comms-adviser-query
# left implicit; if main.ts fatals during runBootCheck, this fires before
# Restart kicks in and gives systemd a clean failed-state. 300 iterations
# matches TimeoutStartSec=300 (one probe/second).
ExecStartPost=/bin/bash -c 'for i in $(seq 1 300); do [ -S /run/ai-calendar-adviser/query.sock ] && exit 0; sleep 1; done; exit 1'

# RuntimeDirectory creates /run/ai-calendar-adviser/ owned by ubuntu:ubuntu
# at service start; systemd cleans it up on stop. main.ts binds its
# socket inside this dir at /run/ai-calendar-adviser/query.sock.
RuntimeDirectory=ai-calendar-adviser
RuntimeDirectoryMode=0755

# Per-user OAuth refresh token (G6.5a + G6.5b) handed in via
# LoadCredential. The runtime path resolves via
# $CREDENTIALS_DIRECTORY/oauth-token. The actual env-var wiring
# (CALENDAR_DB_PATH, CALENDAR_SOCKET_PATH, calendar slot ids,
# OAUTH_TOKEN_PATH) is done in .env. OAUTH_TOKEN_PATH is set in .env to
# $CREDENTIALS_DIRECTORY/oauth-token (resolved at start time). The
# refresh token is bound to ai@liao.info; per
# feedback_no_kelvin_account_impersonation, this agent never
# authenticates as Kelvin.
LoadCredential=oauth-token:/etc/ai-calendar-adviser/oauth-token.json

# Hardening (mirrors ai-comms-adviser).
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
# /var/lib/ai-ops holds the shared control-plane ops.db. ProtectSystem=strict
# makes everything outside ReadWritePaths read-only, which is why the boot-time
# ops.db open returned "attempt to write a readonly database" (AI1) — add it.
ReadWritePaths=/var/lib/ai-calendar-adviser /var/lib/ai-ops
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
# UDS-only — no inbound IP. The agent makes outbound HTTPS to Google
# Calendar for the sync runner (sub-item 3); AF_INET6 + AF_INET cover
# that without permitting raw or netlink sockets.
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
LockPersonality=true
RestrictRealtime=true

[Install]
WantedBy=multi-user.target
