# deploy/systemd/ai-calendar-adviser-pull-deploy.timer
[Unit]
Description=ai-calendar-adviser pull-deploy poller (every 2 min)
Documentation=https://github.com/hamlet-archer/ai-calendar-adviser

[Timer]
OnBootSec=30s
OnUnitActiveSec=2min
AccuracySec=15s
Unit=ai-calendar-adviser-pull-deploy.service

[Install]
WantedBy=timers.target
